# supervisord.conf
[program:trophybot]
command=npm start
directory=%(here)s
autostart=true
autorestart=true
stderr_logfile=logs/trophybot.err.log
stdout_logfile=logs/trophybot.out.log
user=www-data
